# conf/example.conf
[program:example]
autorestart=true
command=/conf/ssh.sh 192.168.1.2 3306
